# ansible/src/calibre.service
[Unit]
Description=Calibre web

[Service]
Type=simple
User=calibre
ExecStart=/usr/local/bin/cps

[Install]
WantedBy=multi-user.target
